# supervisord/cron.conf
[program:cron]
command=/data/scripts/start_cron.sh
autorestart=true
