# supervisord/postgres.conf
[program:postgres]
command=/data/scripts/start_postgres.sh
autorestart=true
