# supervisord/syslog.conf
[program:syslog]
priority=1
command=/usr/sbin/rsyslogd -n
autorestart=true
